# coverage/navigator-layer.json
{
  "name": "splunk-detection-lab coverage",
  "versions": { "attack": "14", "navigator": "5.0.0", "layer": "4.5" },
  "domain": "enterprise-attack",
  "description": "Detection coverage for the splunk-detection-lab repo. 18 detections in testing status. Regenerated from coverage/coverage.md.",
  "filters": { "platforms": ["Windows"] },
  "sorting": 3,
  "layout": { "layout": "side", "showName": true, "showID": false, "showAggregateScores": false, "countUnscored": false },
  "techniques": [
    { $1100, "comment": "PowerShell — encoded + download cradle + AMSI bypass" },
    { "techniqueID": "T1059.003", "score": 50, "comment": "cmd caret obfuscation" },
    { $1100, "comment": "Run/RunOnce key modification" },
    { "techniqueID": "T1053.005", "score": 50, "comment": "Scheduled task via schtasks" },
    { "techniqueID": "T1136.001", "score": 50, "comment": "Local account creation EID 4720" },
    { "techniqueID": "T1546.012", "score": 50, "comment": "IFEO Debugger persistence" },
    { "techniqueID": "T1547.014", "score": 50, "comment": "Active Setup StubPath" },
    { $1100, "comment": "Rundll32 unusual parent / no DLL" },
    { "techniqueID": "T1218.005", "score": 50, "comment": "Mshta remote/inline" },
    { "techniqueID": "T1218.010", "score": 50, "comment": "Regsvr32 Squiblydoo" },
    { "techniqueID": "T1027",     "score": 50, "comment": "Obfuscation entropy on PowerShell" },
    { "techniqueID": "T1140",     "score": 50, "comment": "Certutil decode" },
    { "techniqueID": "T1562.001", "score": 50, "comment": "Defender tampering + AMSI bypass" },
    { "techniqueID": "T1222.001", "score": 50, "comment": "icacls Everyone:F" },
    { "techniqueID": "T1112",     "score": 50, "comment": "Persistence-helper registry keys" },
    { $1100, "comment": "LSASS process access (EID 10)" },
    { "techniqueID": "T1110.001", "score": 50, "comment": "Brute force burst on 4625" },
    { "techniqueID": "T1110.003", "score": 50, "comment": "Password spray" },
    { "techniqueID": "T1105",     "score": 50, "comment": "Curl/bitsadmin + cradle + certutil download" }
  ],
  "gradient": { "colors": ["#ffffff","#ffe066","#3ad55c"], "minValue": 0, "maxValue": 100 },
  "legendItems": [
    { "label": "Production (validated)", "color": "#3ad55c" },
    { "label": "Testing (written, not validated)", "color": "#ffe066" }
  ]
}
